# deploy_artifacts/nvcf_func_monitor_20260701/nvcf_func_monitor.timer
[Unit]
Description=Run NVCF function health monitor every 10 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
